# systemd/media-backup.service
[Unit]
Description=Backup container appdata and sync to OneDrive
After=network-online.target
Wants=network-online.target
RequiresMountsFor=/mnt/media

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/srv/containers/media
ExecStart=/srv/containers/media/scripts/backup-media-appdata.sh
ExecStartPost=/srv/containers/media/scripts/sync-backups-to-onedrive.sh
