# ansible/roles/bookbuddy_worker/templates/bookbuddy-worker.timer.j2
[Unit]
Description=Poll for BookBuddy generation jobs

[Timer]
OnCalendar={{ bookbuddy_worker_interval }}
Persistent=true
RandomizedDelaySec=5
Unit=bookbuddy-worker.service

[Install]
WantedBy=timers.target
